# .config/systemd/user/foot-server.service.d/override.conf
[Service]
EnvironmentFile=-%h/.config/environment.d/environment.conf
